# .docker/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P

exec docker-php-entrypoint php-fpm
